# webservers/weasimov.target
# See ReADME.md for more about this file
[Unit]
Description=Weasimov

[Install]
WantedBy=multi-user.target